# broadlink-ac-mqtt/rootfs/etc/services.d/broadlink-ac-mqtt/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Stop the service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================
